# aycep/finals/verysecretstorage/initramfs/etc/inittab
::sysinit:/etc/init.d/rcS
::once:-sh -c 'clear; setuidgid 1000 sh; poweroff'
